# build/images/centinel/etc/cron.hourly/centinel
#!/bin/bash
# cron job for centinel
/usr/local/bin/centinel --sync
/usr/local/bin/centinel
/usr/local/bin/centinel --sync
