# build/images/centinel/etc/cron.daily/centinel-autoupdate
#!/bin/bash
# autoupdater for centinel
sudo pip install --upgrade centinel
